# ntpdate.upstart
description "Synchronizes time with ntpdate (NTP client)"

# TODO: act on "have default gateway"
# start on net-device-added INTERFACE!=lo
# TODO: run before ntpd

start on pld.network-started
stop on pld.shutdown-started

task

# mark that ntpdate fail is ok (task marked as done and hwclock skipped)
normal exit 0 1

exec /usr/sbin/ntpdate-wrapper
